# deploy/rapid-plateau-dashboard.timer
[Unit]
Description=Weekly schedule for the Rapid Plateau dashboard batch

[Timer]
# Weekly, Sunday 02:00 in the system's local timezone. The batch runs for a few
# hours, so pick a start time that finishes BEFORE any daily maintenance window
# (e.g. a package-update-driven service restart). Confirm the host timezone and
# adjust OnCalendar accordingly.
OnCalendar=Sun *-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
